# scripts/com.agent-pool.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.agent-pool.daemon</string>

	<key>ProgramArguments</key>
	<array>
		<string>AGENT_POOL_BINARY</string>
		<string>start</string>
		<string>POOL_DIR</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<false/>

	<key>StandardOutPath</key>
	<string>POOL_DIR/launchd-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>POOL_DIR/launchd-stderr.log</string>

	<key>WorkingDirectory</key>
	<string>POOL_DIR</string>

	<key>ExitTimeOut</key>
	<integer>35</integer>
</dict>
</plist>
